# amp-control.service
[Unit]
Description=Auto control amplifier when audio sinks are idle (with debounce)
After=network.target sound.target pipewire.service
Requires=pipewire.service 

[Service]
Type=simple
ExecStart=/usr/local/bin/amp-control.sh
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

